# contrib/frozen-bottle-visualizer.service
[Unit]
Description=Frozen Bottle Equalizer

[Service]
ExecStart=/home/labor/devel/FrozenBottle/UDP_clients/audio_visualizer.py frozen-bottle
